# etc/craned.service.in
[Unit]
Description=Crane
After=network.target nss-lookup.target

[Service]
User=root
ExecStart=@BINDIR_VAR@/craned

[Install]
WantedBy=multi-user.target
