# etc/cranectld.service.in
[Unit]
Description=CraneCtld
After=network.target nss-lookup.target

[Service]
User=root
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/cranectld

[Install]
WantedBy=multi-user.target
