# etc/craned.service.in
[Unit]
Description=Crane
After=network.target nss-lookup.target

[Service]
User=root
Group=root
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/craned

[Install]
WantedBy=multi-user.target
